# usb.gadget.service
[Unit]
After=rc-local.service

[Service]
ExecStart=/root/usb.sh
Type=oneshot
RemainAfterExit=true

[Install]
WantedBy=basic.target
